# debian/statsd.upstart
description "statsd"
author      "etsy"

start on startup
stop on shutdown

script
    chdir /opt/statsd

    exec sudo -u nobody /opt/statsd/scripts/start
end script
